# stockbot.service
[Unit]
Description=Indian Stock Market Alert Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/path/to/stock_bot
EnvironmentFile=/path/to/stock_bot/.env
ExecStart=/path/to/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/stock_bot/logs/bot.log
StandardError=append:/path/to/stock_bot/logs/bot.log

[Install]
WantedBy=multi-user.target
